# files/etc/init.d/blue-merle-touch
#!/bin/sh /etc/rc.common
# blue-merle v2 — touchscreen long-press trigger daemon
#
# Runs at START=81, right after gl_screen (S80). blue-merle-touch reads
# /dev/input/event0 non-exclusively (gl_screen keeps working) and forks
# `blue-merle sim-swap` on a 2-second hold of the clock area (X:0-80, Y:0-30).
# Toggled via the LuCI page or: /etc/init.d/blue-merle-touch disable && stop

START=81
STOP=20
USE_PROCD=1

start_service() {
    procd_open_instance
    procd_set_param command /usr/bin/blue-merle-touch
    procd_set_param respawn 30 5 0
    procd_set_param stderr 1
    procd_close_instance
}
